# templates/znc.service
[Unit]
Description=znc

[Service]
Type=forking
ExecStart=/usr/bin/znc
ExecStop=/usr/bin/pkill znc
Restart=always

[Install]
WantedBy=default.target
